# deploy/ditto-screener.service
[Unit]
Description=Ditto SN118 screener worker (queue -> build gate -> verdict)
Documentation=https://github.com/ditto-assistant/ditto-screener
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
User=deploy
Group=ditto
WorkingDirectory=/opt/ditto/screener/src
EnvironmentFile=/opt/ditto/screener/screener.env
Environment=HOME=/opt/ditto/screener
ExecStart=/opt/ditto/screener/src/.venv/bin/ditto-screener
Restart=always
RestartSec=15
StandardOutput=append:/opt/ditto/logs/ditto-screener.log
StandardError=append:/opt/ditto/logs/ditto-screener.log
LimitNOFILE=65536
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
